# linux/src/prBetriebssystem/backup.timer
[Unit]
Description=Back-up für Home-Verzeichnis 

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h
Unit=backup.service

[Install]
WantedBy=timers.target
