# linux/src/prBetriebssystem/backup.service
[Unit]
Description= Ausführen von backup.sh

[Service]
Type=simple
ExecStart=/etc/backup.sh

[Install]
WantedBy=multi-user.target
